# go_pinger/supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=sh -c "cd /app/api && ./api"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:request]
command=sh -c "cd /app/request && ./request"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
